# live-build/config/includes.chroot/etc/systemd/system/ollamalinux-monitor.service
[Unit]
Description=OllamaLinux System Health Monitor
After=ollama.service

[Service]
Type=simple
ExecStart=/usr/local/bin/ollamalinux-monitor live
Restart=on-failure
RestartSec=10
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ReadOnlyPaths=/
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
