# live-build/config/includes.chroot/etc/systemd/system/ollama.service
[Unit]
Description=Ollama LLM Server
Documentation=https://ollama.com
After=network-online.target
Wants=network-online.target
ConditionPathExists=/etc/ollamalinux/.firstboot-done

[Service]
Type=simple
User=ollama
Group=ollama
EnvironmentFile=/etc/default/ollama
ExecStart=/usr/local/bin/ollama serve
Restart=always
RestartSec=3
LimitNOFILE=1048576
LimitMEMLOCK=infinity
LimitNPROC=4096
Nice=-5
IOSchedulingClass=realtime
IOSchedulingPriority=0
NoNewPrivileges=true
SupplementaryGroups=video render
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/ollama /tmp

[Install]
WantedBy=multi-user.target
