# live-build/config/includes.chroot/etc/systemd/system/open-webui.service
[Unit]
Description=Open WebUI - AI Chat Interface
Documentation=https://github.com/open-webui/open-webui
After=network-online.target ollama.service
Wants=ollama.service
Requires=network-online.target
ConditionPathExists=/etc/ollamalinux/.firstboot-done

[Service]
Type=simple
User=openwebui
Group=openwebui
EnvironmentFile=/etc/ollamalinux/open-webui.conf
WorkingDirectory=/opt/open-webui
ExecStart=/opt/open-webui/venv/bin/open-webui serve --host 127.0.0.1 --port 8080
Restart=always
RestartSec=5
LimitNOFILE=65536
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
Environment=OLLAMA_BASE_URL=http://127.0.0.1:11434
Environment=WEBUI_AUTH=true
Environment=DATA_DIR=/opt/open-webui/data

[Install]
WantedBy=multi-user.target
